# data/supervisor_django.conf
# 本文件用于supervisor 配置守护uwsgi进程
# 本文件放置在/etc/supervisor/conf.d/目录下

[program:django]
command = python manage.py runserver 0.0.0.0:8000
directory = /home/workspace
autostart=true
autorestart=true
stdout_logfile=/var/log/django.log
redirect_stderr=true
stopsignal=QUIT